# scripts/smart-box-ad-sync.service
[Unit]
Description=Smart Box Active Directory Sync
After=smart-box.service docker.service network-online.target
Wants=smart-box.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=__SMART_BOX_USER__
WorkingDirectory=__SMART_BOX_PROJECT_ROOT__
ExecStart=/usr/bin/env bash __SMART_BOX_PROJECT_ROOT__/scripts/sync_ad_users.sh
