# scripts/smart-box-reminder.timer
[Unit]
Description=Daily Smart Box Return Reminder at 18:00

[Timer]
OnCalendar=*-*-* 18:00:00
Persistent=true
Unit=smart-box-reminder.service

[Install]
WantedBy=timers.target
