# scripts/smart-box-rc522-reader.service
[Unit]
Description=Smart Box RC522 RFID Reader
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=__SMART_BOX_USER__
WorkingDirectory=__SMART_BOX_PROJECT_ROOT__
ExecStart=/usr/bin/env bash __SMART_BOX_PROJECT_ROOT__/scripts/start-rc522-reader.sh
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
